# plans/systemd/myaudible.service
[Unit]
Description=myAudible Processing Service
After=network.target myaudible-input.path
Requires=myaudible-input.path

[Service]
Type=simple
User=myaudible
Group=myaudible
WorkingDirectory=/opt/myaudible
ExecStart=/opt/myaudible/venv/bin/python -m myaudible.app
EnvironmentFile=/opt/myaudible/.env
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
# Logging configuration
ExecStartPre=/bin/mkdir -p /opt/myaudible/logs

[Install]
WantedBy=multi-user.target
