# plans/systemd/myaudible.service
[Unit]
Description=Private Reading Processing Service
After=network.target private-reading-input.path
Requires=private-reading-input.path

[Service]
Type=simple
User=private-reading
Group=private-reading
WorkingDirectory=/opt/private-reading
ExecStart=/opt/private-reading/venv/bin/python -m private_reading.app
EnvironmentFile=/opt/private-reading/.env
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
# Logging configuration
ExecStartPre=/bin/mkdir -p /opt/private-reading/logs

[Install]
WantedBy=multi-user.target
